# backend/daemon/cmasd.service
[Unit]
Description=OpenCMAS service 
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/bin/node /home/stefan/opencmas5/opencmas/backend/daemon/main.js

[Install]
WantedBy=multi-user.target1