# backend/daemon/cmasd.service
[Unit]
Description=OpenCMAS service 

[Service]
User=root
RemainAfterExit=yes
WorkingDirectory=/opt/opencmas/Frontend
ExecStart=/bin/bash /opt/opencmas/Frontend/main.sh

[Install]
WantedBy=multi-user.target


