# backend/daemon/cmasd.service
[Unit]
Description=OpenCMAS service 

[Service]
User=root
RemainAfterExit=yes
WorkingDirectory=/bin/opencmas
ExecStart=/bin/bash /bin/opencmas/backend/daemon/main.sh

[Install]
WantedBy=multi-user.target


